# systemd/system/zram_swap.service
[Unit]
Description=swap with zram
DefaultDependencies=no
Before=swap.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/zram-init -d0 -- 512
ExecStop=/usr/sbin/zram-init -d0 -- 0

[Install]
WantedBy=swap.target
